# infra/ansible/roles/common/templates/node_exporter.service.j2
{#
  node_exporter systemd unit.
  - Runs as an unprivileged system user (ProtectHome/ProtectSystem
    hardening is meaningful because of this, not in spite of it).
  - Only the `systemd` collector is disabled by default on
    node_exporter but we disable a few more to keep the scrape
    response small on a 1GB box.
  - Binds to :9100 on all interfaces — the SG on the monitoring
    tier restricts inbound to the prometheus host.
#}
[Unit]
Description=Prometheus Node Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network-online.target
Wants=network-online.target

[Service]
User={{ common_node_exporter_user }}
Group={{ common_node_exporter_user }}
Type=simple
ExecStart={{ common_node_exporter_bin }} \
  --web.listen-address=:{{ node_exporter_port }} \
  --collector.disable-defaults \
  --collector.cpu \
  --collector.meminfo \
  --collector.diskstats \
  --collector.filesystem \
  --collector.loadavg \
  --collector.netdev \
  --collector.systemd \
  --collector.uname \
  --collector.vmstat
Restart=on-failure
RestartSec=5s

# Hardening — reasonable defaults for a metrics exporter that only
# needs to read /proc and /sys.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
